# server/etc/systemd/system/izipublish_sync.service
[Unit]
Description=Easy Publish Update Chain Sync Trigger
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/izipublish_sync
ExecStart=/usr/bin/env bash -lc "${IZIPUBLISH_SYNC_SCRIPT:-/root/izipublish/cron/izipublish_sync.sh}"
