# server/etc/systemd/system/izipublish_sync.service
[Unit]
Description=Sync Update Chain Script

[Service]
Type=oneshot
ExecStart=sh /root/izipublish/cron/izipublish_sync.sh
